# farnsworth.cron
*/5 * * * * www-data /opt/apps/farnsworth/manage.py runcrons
