# dockerfiles/farnsworth.cron
*/5 * * * * root /opt/apps/farnsworth/manage.py runcrons
